# scripts/com.brdavies.everybody-llmbo.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.brdavies.everybody-llmbo</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/brdavies/.codex/bin/start-codex-bridge.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/brdavies/Library/Logs/everybody-llmbo.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/brdavies/Library/Logs/everybody-llmbo.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
